# gitgrope.service
[Unit]
Description=Grope it like it's hot

[Service]
ExecStart=/usr/sbin/gitgrope

[Install]
WantedBy=multi-user.target